# deploy/edge/hookly-edge.service
[Unit]
Description=Hookly Edge Gateway
After=docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/hookly-edge
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
